# dot_config/systemd/user/opencode-web.service
[Unit]
Description=OpenCode Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
Environment=HOME=%h
Environment=XDG_CONFIG_HOME=%h/.config
ExecStart=%h/.opencode/bin/opencode web --hostname 0.0.0.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
